# goals/supervisord.conf
[program:runserver]
command={{ PYTHON }} {{ PROJECT_DIR }}/manage.py runserver

[program:redis]
command=redis-server redis.conf
directory={{ PROJECT_DIR }}/../
